# contrib/ctest-templates/test-umundo-experimental.cron
#!/bin/bash
export CTEST_SUBMIT_TYPE="Experimental"

if [ ! -n "${UMUNDO_SOURCE_DIR:+x}" ]; then
	export UMUNDO_SOURCE_DIR="/var/builds/umundo"
fi

ME=`basename $0`
SCRIPT_DIR="$( cd "$( dirname "$0" )" && pwd )"
LOCK_DIR=/tmp/umundo-build.lock

# try to find ctest
CTEST=/usr/local/bin/ctest

cd $SCRIPT_DIR
TESTS=`ls umundo*.ctest`
while true; do
	if mkdir "$LOCK_DIR"; then
		trap 'rm -rf "$LOCK_DIR"' 0    # remove directory when script finishes
		for TEST in ${TESTS}; do
			nice ${CTEST} -VV --timeout 60 -S ${TEST}
		done
		exit
	else
		echo "Could not get build lock- waiting"
		sleep 10
	fi
done
